# stacks/daemon/openclaw-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Adjust path to your OpenClaw installation -->
        <string>/usr/local/bin/openclaw</string>
        <string>gateway</string>
        <string>start</string>
    </array>

    <!-- Start on load and restart on crash -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Restart throttle: wait 10s between restarts -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>OPENCLAW_GATEWAY_TOKEN</key>
        <string>YOUR_TOKEN_HERE</string>
        <!-- Add other env vars as needed -->
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/openclaw-gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/openclaw-gateway-error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USER/.openclaw</string>

    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
